# in-target/packages/files/.config/autostart/stretchly.desktop
[Desktop Entry]
Type=Application
Name=stretchly
Exec=/usr/bin/stretchly
StartupNotify=false
Terminal=false
