# configs/anacron/anacrontab
# /etc/anacrontab: configuration file for anacron

# See anacron(8) and anacrontab(5) for details.

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
HOME=/root
LOGNAME=root
MAILTO=zubmic
RANDOM_DELAY=30

# These replace cron's entries
1               5       cron.daily      run-parts --report /etc/cron.daily
7               10      cron.weekly     run-parts --report /etc/cron.weekly
@monthly        15      cron.monthly    run-parts --report /etc/cron.monthly

# User defined job
1               5       clam.daily      /usr/bin/clamscan -ir /home
7               5       clam.weekly     /usr/bin/clamscan -ir /